# docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/storage/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=1
stdout_logfile=/app/storage/logs/php-fpm.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/app/storage/logs/php-fpm-error.log
stderr_logfile_maxbytes=5MB

[program:queue-worker]
command=php /app/artisan queue:work --queue=default --sleep=3 --tries=2 --max-jobs=500
directory=/app
autostart=true
autorestart=true
; Restart after 500 jobs to prevent memory leaks (supervisor will relaunch automatically)
startsecs=5
stopwaitsecs=30
priority=10
stdout_logfile=/app/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/app/storage/logs/queue-worker-error.log
stderr_logfile_maxbytes=5MB
